# stage2-1-beepy-config/00-optimize-boot/files/load-brcmfmac.service
[Unit]
Description=Load brcmfmac Kernel Module
DefaultDependencies=no
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/sbin/modprobe brcmfmac

[Install]
WantedBy=post-boot.target
