# stage2-1-beepy-config/00-optimize-boot/files/post-boot.target
[Unit]
Description=Start background services after multi-user
Requires=multi-user.target
After=multi-user.target
