# stage2-beepy/01-optimize-boot/files/defer-systemd-timesyncd.service
[Unit]
Description=Load systemd-timesyncd
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl start systemd-timesyncd

[Install]
WantedBy=post-boot.target
